# ec2/roles/server/templates/music.mount.j2
[Unit]
Description=Mount music volume

[Mount]
What=/dev/xvdf
Where={{ mount_point }}
Type={{ build_fs_type }}
Options=defaults
