# ec2/roles/worker/templates/music-snapshot.service.j2
[Unit]
Description=Snapshot build volume

[Service]
Type=oneshot
ExecStart=/usr/bin/aws ec2 create-snapshot
