# ec2/roles/worker/templates/music-sources.mount.j2
[Unit]
Description=Mount music volume

[Mount]
What={{ sources_device }}
Where={{ sources_mount_point }}
Type={{ sources_fs_type }}
Options=defaults
