# ec2/roles/worker/templates/music-builds.mount.j2
[Unit]
Description=Mount music volume

[Mount]
What={{ builds_device }}
Where={{ builds_mount_point }}
Type={{ builds_fs_type }}
Options=defaults
