# ec2/roles/worker/templates/music-monitor.service.j2
[Unit]
Description=Monitor music files for changes

[Service]
Type=simple
ExecStart=/usr/bin/tup monitor
User=music
Group=music
WorkingDirectory={{ work_dir }}
Restart=on-failure
